# scripts/templates/kaleido-api-nginx.container.jinja2
[Unit]
Description=nginx behind Cloudflare Tunnel
Requires=kaleido-api.network
Requires=kaleido-api-content.volume

[Container]
Image=docker.io/library/nginx:alpine
ContainerName=kaleido-api-nginx
Network=kaleido-api.network:alias=nginx

Volume={{ volume_dir }}:/etc/nginx/conf.d:ro
Volume=kaleido-api-content.volume:/usr/share/nginx/html:ro

NoNewPrivileges=true
ReadOnly=true
Tmpfs=/var/cache/nginx
Tmpfs=/var/run
Tmpfs=/tmp

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target